# vizitka-bot.service
[Unit]
Description=Vizitka Qualification Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/projects/vizitka-bot
EnvironmentFile=/root/projects/vizitka-bot/.env
ExecStart=/root/projects/vizitka-bot/.venv/bin/python bot.py
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
